# templates/kea-dhcp4.service.j2
[Unit]
Description=Kea DHCPv4 Server
Wants=network-online.target
After=network-online.target

[Service]
ExecStart={{kea_path}}/sbin/kea-dhcp4 -c {{kea_config}}

[Install]
WantedBy=multi-user.target
